# supervisord.conf
[supervisord]
loglevel=debug
logfile=/dev/stdout
logfile_maxbytes=0

[inet_http_server]
port=127.0.0.1:9001

[supervisord]
loglevel=debug
logfile=/dev/stdout
logfile_maxbytes=0

[inet_http_server]
port=127.0.0.1:9001

# attestation server
[program:attestation-server]
command=/app/attestation-server --ip-addr 127.0.0.1:1300 --pub-key /app/id.pub
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# attestation server proxy
[program:attestation-proxy]
command=/app/vsock-to-ip --vsock-addr 88:1300 --ip-addr 127.0.0.1:1300
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# attestation server
[program:oyster-attestation-utility]
command=/app/attestation-server --ip-addr 127.0.0.1:1500 --pub-key /app/secp.pub
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# attestation server proxy
[program:oyster-attestation-utility-proxy]
command=/app/vsock-to-ip --vsock-addr 88:1500 --ip-addr 127.0.0.1:1500
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# transparent proxy component inside enclave
[program:ip-to-vsock-transparent]
command=/app/ip-to-vsock-transparent --vsock-addr 3:1200 --ip-addr 127.0.0.1:1200
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# DNS-over-HTTPS provider
[program:dnsproxy]
command=/app/dnsproxy -u https://1.1.1.1/dns-query -v
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0c

# #Genertor client
# [program:generator_client]
# command=/app/generator_client
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stdout
# stderr_logfile_maxbytes=0
# environment=PORT=5000,SUPERVISORD_PATH=./app/supervisord,GENERATOR_PATH=zkbob_generator,SKIP_VERIFICATION=true

# # transparent proxy component for generator client
# [program:generator-client-proxy]
# command=/app/vsock-to-ip --vsock-addr 88:5000 --ip-addr 127.0.0.1:5000
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stdout
# stderr_logfile_maxbytes=0

# risc0 prover
[program:risc0-generator]
command=/app/kalypso-attestation-prover
autorestart=true
autostart=true
environment=PORT=3030
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
environment=PRIVATE_KEY=6825e7cc45c307480f5fea2aa9ea35a4052e6a4db5fcd062b9f043fa924dfb61,RPC_URL=https://arb-sepolia.g.alchemy.com/v2/cFwacd_RbVpNrezyxZEvO6AnnCuO-kxt,CHAIN_ID=421614,GENERATOR_REGISTRY_ADDRESS="0x5ce3e1010028C4F5687356D721e3e2B6DcEA7C25",PROOF_MARKETPLACE_ADDRESS="0x0b6340a893B944BDc3B4F012e934b724c83abF97",REWARD_ADDRESS=0x9674f32Be56c467594a1B36dEB1E35D513EA0957,MARKET_ID=26,START_BLOCK=95285467,MAX_PARALLEL_PROOFS=1,IVS_URL=http://3.110.146.109:3030,PROVER_URL=http://localhost:3030/api/generateProof,GAS_KEY=6825e7cc45c307480f5fea2aa9ea35a4052e6a4db5fcd062b9f043fa924dfb61

# transparent proxy component for zkbob prover
[program:risc0-prover-proxy]
command=/app/vsock-to-ip --vsock-addr 88:3030 --ip-addr 127.0.0.1:3030
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# transparent proxy component for health check
[program:risc0-health-check-proxy]
command=/app/vsock-to-ip --vsock-addr 88:9999 --ip-addr 127.0.0.1:9999
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
